# cron/crontab
* * * * * curl -s http://web/check_for_sessions_to_stop.php >> /var/log/cron.log 2>&1
